# hack/docker/s6/backpulse-serve/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

export USER=backpulse
exec gosu $USER /app/backpulse/backpulse
