# hack/docker/s6/.s6-svscan/finish
#!/bin/sh

# Cleanup some services and s6 event folder
